# systemd/motion_wide.service
[Unit]
Description=Motion_No_Cam Flask Service
After=network-online.target i2c-dev.service
Wants=network-online.target i2c-dev.service

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/opt/Motion_No_Cam
Environment=FLASK_ENV=production
Environment=PYTHONUNBUFFERED=1
Environment=MOTION_DB=/opt/Motion_No_Cam/motion.db

# Run migrations as root (same user as service)
ExecStartPre=/opt/Motion_No_Cam/venv/bin/python -m app.migrations
# Wait for /dev/i2c-1 to appear
ExecStartPre=/bin/sh -c 'for i in $(seq 1 10); do [ -e /dev/i2c-1 ] && exit 0; sleep 1; done; exit 1'

ExecStart=/opt/Motion_No_Cam/venv/bin/python -m app.main

Restart=on-failure
RestartSec=2

# Keep FS mostly read-only but allow app path
ProtectSystem=full
ReadWritePaths=/opt/Motion_No_Cam

[Install]
WantedBy=multi-user.target
